# contrib/systemd/musefs-scan.timer
[Unit]
Description=Periodic musefs library re-scan

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
