# contrib/systemd/musefs.service
[Unit]
Description=musefs read-only re-tagging FUSE mount
Documentation=https://github.com/Sohex/musefs

[Service]
Type=simple
EnvironmentFile=-%h/.config/musefs/musefs.conf
# The --user manager does not inherit a login shell's PATH, so a cargo-installed
# binary in ~/.cargo/bin is not found by a bare `musefs`. Adjust this PATH (or
# replace ExecStart with an absolute path) to match where musefs is installed.
Environment=PATH=%h/.cargo/bin:/usr/local/bin:/usr/bin
ExecStart=musefs mount
# Passing --template here? systemd expands `$word` in ExecStart as a (usually
# undefined, hence empty) environment variable, so a bare `$albumartist/...`
# reaches musefs garbled and the mount aborts (rc 2). Double the `$` so systemd
# passes the literal template through:
#ExecStart=musefs mount %h/mnt/music --db %h/.local/share/musefs/library.db --template '$$albumartist/$$album/$$title'
# musefs unmounts cleanly on SIGTERM (systemd's default stop signal), so no
# ExecStop is required. Uncomment as a fallback if a mount is ever left behind:
#ExecStop=-fusermount3 -u ${MUSEFS_MOUNTPOINT}
Restart=on-failure
RestartSec=5
# This unit is deliberately NOT sandboxed, and that cannot be changed without
# breaking the mount. musefs mounts via the *setuid* fusermount3 helper, and
# NoNewPrivileges=true disables the setuid escalation, so the mount fails with
# "fusermount3: mount failed: Operation not permitted". Crucially, nearly every
# systemd hardening directive (SystemCallFilter=, Restrict*, LockPersonality=,
# MemoryDenyWriteExecute=, ProtectKernel*, ProtectControlGroups=, ...) forces
# the kernel no_new_privs flag — it is required to install a seccomp filter for
# an unprivileged process — which breaks the mount the same way, *even with*
# NoNewPrivileges=no. So none of them can be added here.
# Do NOT add ProtectHome=, PrivateMounts=, or MountFlags=private either: those
# place the mount in a private namespace and hide it from the rest of your
# session.
# (The scanner, musefs-scan.service, creates no mount and IS fully sandboxed.)

[Install]
WantedBy=default.target
